# modal/cpu_pipeline/supervisord.conf
; Topologi C — kör alla fem tjänster i en container. Varje gunicorn körs i SIN
; egen cwd för att undvika att tjänsternas `main`/`src`-moduler krockar.
; Portar matchar SERVICE_URL:erna i app.py (localhost).
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
loglevel=info

; --- CPU-tjänster (interna) ---------------------------------------------------
[program:encoder]
command=gunicorn --bind :8082 --workers 1 --threads 8 --timeout 900 src.main:app
directory=/app/encoder
environment=PYTHONPATH="/app/encoder",PORT="8082"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:obstruction]
command=gunicorn --bind :8081 --workers 2 --threads 2 --worker-class gthread --timeout 900 main:app
directory=/app/obstruction/src
environment=PYTHONPATH="/app/obstruction/src",PORT="8081"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:merger]
command=gunicorn --bind :8084 --workers 1 --threads 8 --timeout 900 main:app
directory=/app/merger/src
environment=PYTHONPATH="/app/merger/src",PORT="8084"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:stats]
command=gunicorn --bind :8085 --workers 1 --threads 8 --timeout 900 main:app
directory=/app/stats/src
environment=PYTHONPATH="/app/stats/src",PORT="8085"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; --- Gateway (publik via Modal web_server på 8080) ----------------------------
[program:lux]
command=gunicorn --bind :8080 --workers 1 --threads 8 --timeout 900 main:app
directory=/app/lux/src
environment=PYTHONPATH="/app/lux/src",PORT="8080"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
